# deploy/systemd/containqr.service
[Unit]
Description=ContainQR gunicorn service
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/srv/containqr/app
EnvironmentFile=/srv/containqr/env/.env
ExecStart=/srv/containqr/venv/bin/gunicorn backend.wsgi:application \
  --bind 127.0.0.1:8001 \
  --workers 3 \
  --timeout 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
